# deploy/systemd/okoznaniy-monitor.timer
[Unit]
Description=Run OkoZnaniy health monitor every 2 minutes

[Timer]
OnBootSec=45s
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true
Unit=okoznaniy-monitor.service

[Install]
WantedBy=timers.target
